# setup/smartos/resources/jenkins_manifest.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='jenkins'>
    <service name='application/jenkins' type='service' version='1'>

        <create_default_instance enabled='true' />

        <single_instance />

        <dependency name='network' grouping='require_all' restart_on='error' type='service'>
            <service_fmri value='svc:/milestone/network:default' />
        </dependency>
        <dependency name='filesystem' grouping='require_all' restart_on='error' type='service'>
            <service_fmri value='svc:/system/filesystem/local' />
        </dependency>

        <method_context working_directory="/home/iojs">
            <method_credential user='iojs' group='iojs' privileges='basic' />
            <method_environment>
                <envvar name='NODE_COMMON_PIPE' value='/home/iojs/test.pipe' />
                <envvar name='OSTYPE' value='solaris' />
                <envvar name='HOME' value='/home/iojs' />
                <envvar name='JOBS' value='8' />
                <envvar name='DESTCPU' value='{{destcpu}}' />
                <envvar name='PATH' value='/usr/local/libexec/ccache:/usr/local/sbin:/usr/local/bin:/opt/local/sbin:/opt/local/bin:/usr/sbin:/usr/bin:/sbin' />
            </method_environment>
        </method_context>

        <exec_method type='method'
                     name='start'
                     exec='/opt/local/bin/java -Xmx512m -jar slave.jar -jnlpUrl https://jenkins-iojs.nodesource.com/computer/iojs-joyent-smartos-{{id}}/slave-agent.jnlp -secret {{secret}} >> /home/iojs/jenkins_console.log'
                     timeout_seconds='60' />

        <exec_method type="method"
                     name="stop"
                     exec=":kill"
                     timeout_seconds="60" />

        <exec_method name='refresh'
                     type='method'
                     exec=':kill -HUP'
                     timeout_seconds='60' />

        <property_group name="startd" type="framework">
            <propval name="duration" type="astring" value="child"/>
            <propval name="ignore_error" type="astring" value="core,signal"/>
        </property_group>

        <property_group name="application" type="application"></property_group>

        <stability value='Evolving' />

        <template>
            <common_name>
                <loctext xml:lang='C'>Jenkins Continuous Build Server</loctext>
            </common_name>
            <documentation>
                <doc_link name='hudson.dev.java.net' uri='http://jenkins-ci.org/' />
            </documentation>
        </template>
    </service>
</service_bundle>
