# setup/www/resources/config/github-webhook.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

setuid nodejs

exec github-webhook --config /etc/github-webhook.json
